# packages/api/util/ops/spicygolf.cron
# # backups every night, 5am Eastern
 0  5  *  *  *    golf   /usr/bin/nice /home/golf/ops/backup.sh
